# samples/com.github.fxg42.tmux-weather.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.github.fxg42.tmux-weather</string>

    <key>ProgramArguments</key>
    <array>
      <string>/absolute/path/to/clone/path/tmux-weather/scripts/write_weather.sh</string>
      <string>/absolute/path/to/python3</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
